# app/supervisord.conf
[program:node]
command = node launcher.js
directory = /home/dotcloud/current
